# wrt_core/recipes/default_theme_fluent/files/etc/uci-defaults/990_set_fluent_default
#!/bin/sh

uci batch <<'EOF'
set luci.main.mediaurlbase='/luci-static/fluent'
set luci.themes.Fluent='/luci-static/fluent'
commit luci
EOF
